# AnyConnect/build/dnsmasq.service
[Unit]
Description=dnsmasq
Documentation=man:dnsmasq(8)
After=network-online.target
After=dbus.service

[Service]
Type=simple
PrivateTmp=true
ExecStart=/usr/sbin/dnsmasq -d
ExecStop=/bin/kill -KILL $MAINPID
RestartSec=3s
Restart=always

[Install]
WantedBy=multi-user.target
