# AnyConnect/build/ocserv.service
[Unit]
Description=OpenConnect SSL VPN server
Documentation=man:ocserv(8)
After=dbus.service local-fs.target network.target

[Service]
Type=simple
PrivateTmp=true
ExecStartPre=/bin/bash /etc/ocserv/ctl.sh
ExecStart=/usr/sbin/ocserv --foreground --config /etc/ocserv/ocserv.conf
ExecStartPost=/bin/bash /etc/ocserv/ctl.sh CHECK
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -KILL $MAINPID
RestartSec=3s
Restart=always

[Install]
WantedBy=multi-user.target
